# deploy/storyworlds.service
[Unit]
Description=StoryWorlds MUD Server
After=network.target

[Service]
Type=simple
User=storyworlds
Group=storyworlds
WorkingDirectory=/opt/storyworlds
ExecStart=/opt/storyworlds/storyworlds-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=CONFIG_PATH=/opt/storyworlds/config.yaml

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/storyworlds

# Allow binding to port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
